# supervisord_ldap_contacts.ini
[program:laravel-worker]
directory=/var/www/html/ldap_contacts   ; path to site webroot
command=php artisan queue:work
autostart=true
autorestart=true
stderr_logfile=/var/log/ldap_contacts.err.log   ; error log location, need to be exist
stdout_logfile=/var/log/ldap_contacts.app.log   ; app log location, need to be exist
stopsignal=INT
user=root
